# timers/rotateTODO.timer
[Unit]
Description=Update TODO_*.chk file daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
